# templates/systemd.service.template
[Unit]
Description=Node Monitor Backend
After=network.target

[Service]
Type=simple
User={{SYSTEM_USER}}
WorkingDirectory={{APP_DIR}}/backend
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
EnvironmentFile={{APP_DIR}}/backend/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_NAME}}

[Install]
WantedBy=multi-user.target
